# data/scraper_spiegel_online.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>scraperspiegelonline</string>
        <key>ProgramArguments</key>
        <array>
                <string>/Users/munzerts/Dropbox/Uni/Projects/2017-CSS-Attention/code/scraping_pages/scraper_spiegel_online.R</string>
        </array>
	    <key>StartInterval</key>
	    <integer>10000</integer>
</dict>
</plist>