# data/spiegelheadlines.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>spiegelheadlines</string>
        <key>ProgramArguments</key>
        <array>
                <string>/Users/munzerts/github/rscraping-hu-2017/spiegel_scraper.R</string>
        </array>
	    <key>StartInterval</key>
	    <integer>10</integer>
</dict>
</plist>